# packages/avado/build/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:hoprd]
directory=/app/hoprnet/packages/hoprd
command=node ./lib/main.cjs
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
environment=
  DEBUG="hopr*",
  NODE_OPTIONS="--max_old_space_size=4096 --experimental-wasm-modules",
  NODE_ENV="production",
  NEXT_TELEMETRY_DISABLED="1",
  HOPRD_ENVIRONMENT="%(ENV_HOPRD_ENVIRONMENT)s",
  HOPRD_HOST="%(ENV_HOPRD_HOST)s",
  HOPRD_ANNOUNCE="%(ENV_HOPRD_ANNOUNCE)s",
  HOPRD_API="%(ENV_HOPRD_API)s",
  HOPRD_API_HOST="%(ENV_HOPRD_API_HOST)s",
  HOPRD_API_PORT="%(ENV_HOPRD_API_PORT)s",
  HOPRD_HEALTH_CHECK="%(ENV_HOPRD_HEALTH_CHECK)s",
  HOPRD_HEALTH_CHECK_HOST="%(ENV_HOPRD_HEALTH_CHECK_HOST)s",
  HOPRD_HEALTH_CHECK_PORT="%(ENV_HOPRD_HEALTH_CHECK_PORT)s",
  HOPRD_PASSWORD="%(ENV_HOPRD_PASSWORD)s",
  HOPRD_DISABLE_API_AUTHENTICATION="%(ENV_HOPRD_DISABLE_API_AUTHENTICATION)s",
  HOPRD_API_TOKEN="%(ENV_HOPRD_API_TOKEN)s",
  HOPRD_IDENTITY="%(ENV_HOPRD_IDENTITY)s",
  HOPRD_DATA="%(ENV_HOPRD_DATA)s",
  HOPRD_INIT="%(ENV_HOPRD_INIT)s",
  HOPRD_ALLOW_PRIVATE_NODE_CONNECTIONS="%(ENV_HOPRD_ALLOW_PRIVATE_NODE_CONNECTIONS)s",
  HOPRD_PROVIDER="%(ENV_HOPRD_PROVIDER)s"

[program:hopr-admin]
directory=/hopr-admin
command=node server.js
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
environment=
  NODE_ENV="production",
  NEXT_TELEMETRY_DISABLED="1"

[program:prometheus]
command=/opt/prometheus/prometheus
  --web.listen-address=0.0.0.0:9090
  --config.file=/etc/prometheus/prometheus.yml
  --storage.tsdb.path=/opt/prometheus/data
  --web.console.libraries=/opt/prometheus/console_libraries
  --web.console.templates=/opt/prometheus/consoles
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true

[program:grafana]
command=/opt/grafana/bin/grafana-server --config /etc/grafana/grafana.ini --homepath /opt/grafana web
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true


