# deploy/nfpm/hoprd.initd
#!/sbin/openrc-run

description="HOPR node service"
command="/usr/bin/hoprd"
command_args="--configurationFilePath /etc/hoprd/hoprd.cfg.yaml"
pidfile="/var/run/hoprd.pid"
# Source the environment file
if [ -f /etc/hoprd/hoprd.env ]; then
    . <(grep -v '^#' /etc/hoprd/hoprd.env)
    export $(grep -v '^#' /etc/hoprd/hoprd.env | cut -d= -f1) # Export all variables
fi
depend() {
    need net
}
